# doc/config/supervisor.ini
[program:bi_web]
command=<PYTHON_PATH>/python <project_path>/runserver.py
directory=<project_path>
stdout_logfile=<log_path>/app.log
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
stopsignal=KILL
user=root